# deployment/ocrbox.service
[Unit]
Description=OCRBox - Self-Hosted OCR Service
Documentation=https://github.com/stefanahman/ocrbox
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/ocrbox

# Start the service
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up

# Stop the service
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536

# Security
PrivateTmp=yes
NoNewPrivileges=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ocrbox

[Install]
WantedBy=multi-user.target

